# TidyMacHelper/com.tidymac.TidyMacHelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tidymac.TidyMacHelper</string>
    <!-- BundleProgram is the SMAppService-friendly way to point at a
         binary embedded inside the parent app bundle. macOS 13+ resolves
         this relative to the host app, so we don't have to hard-code
         /Applications/TidyMac.app or worry about install paths. -->
    <key>BundleProgram</key>
    <string>Contents/MacOS/com.tidymac.TidyMacHelper</string>
    <key>MachServices</key>
    <dict>
        <key>com.tidymac.TidyMacHelper</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.tidymac.TidyMac</string>
    </array>
</dict>
</plist>
